# dockerfiles/base_disco/supervisord.conf
[supervisord]
nodaemon=true

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface


[program:sshd]
command=/usr/sbin/sshd -D

[program:disco]
command=/disco/bin/disco nodaemon -v
user=disco
environment=HOME="/home/disco"
