# dockerfiles/disco_test/supervisord.conf
[supervisord]
nodaemon=True
loglevel=debug

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface


[program:sshd]
command=/usr/sbin/sshd -D

[program:disco]
command=/disco/bin/disco nodaemon -v
user=disco
environment=HOME="/home/disco"

[program:test]
command=python launch.py run_disco_test.py
directory=/home/mltsp/mltsp/mltsp/docker_scripts
autorestart=unexpected
